# docker/supervisor/php-fpm.conf
[program:php-fpm]
command=/usr/sbin/php-fpm7.3
autorestart=true
