# docker/supervisor/redis.conf
[program:redis]
command=/usr/bin/redis-server
autorestart=true
